# Small_5/package/ipk/luci-app-adblock-plus/root/etc/uci-defaults/luci-adblock
uci -q get firewall.adblock >/dev/null || uci batch <<EOF
	set firewall.adblock=include
	set firewall.adblock.type=script
	set firewall.adblock.path=/var/etc/adblock.include
	commit firewall
EOF
