# Small_5/package/ipk/luci-app-overwall/root/etc/uci-defaults/luci-overwall
uci batch <<EOF
	delete firewall.overwall
	set firewall.overwall=include
	set firewall.overwall.type=script
	set firewall.overwall.path=/var/etc/overwall.include
	set firewall.overwall.reload=1
	commit firewall
EOF
