# Small_5/package/ipk/luci-app-overwall/root/etc/uci-defaults/luci-overwall
uci -q get firewall.overwall >/dev/null || uci batch <<EOF
	set firewall.overwall=include
	set firewall.overwall.type=script
	set firewall.overwall.path=/var/etc/overwall.include
	commit firewall
EOF
