# Small_5/package/ipk/luci-app-https-dns-proxy-whisky/root/etc/hotplug.d/iface/22-https-dns-proxy
#!/bin/sh
if [ "$ACTION" = ifup -a "$INTERFACE" != wan_6 ];then
	[ "$(uci -q get https-dns-proxy.@base[0].enabled)" = 1 ] && service https-dns-proxy restart
fi
